# JobOS/deploy/prisma-web.service
[Unit]
Description=PrismaOS Web UI (Uvicorn)
After=network.target

[Service]
User=szmyt
Group=szmyt
WorkingDirectory=/home/szmyt/server-services/prismaOS
Environment=PYTHONPATH=/home/szmyt/server-services
EnvironmentFile=/home/szmyt/server-services/prismaOS/.env
ExecStart=/home/szmyt/server-services/prismaOS/venv/bin/uvicorn web_ui.app:app --host 0.0.0.0 --port 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
